# meta-enclustra-socfpga/recipes-apps/test-init/files/test.service
[Unit]
Description=Startup script for production test

[Service]
ExecStart=/usr/bin/pt-autostart.sh

[Install]
WantedBy=multi-user.target
